# etc/sv/dbus/run
#!/bin/sh
exec 2>&1
# Create dbus runtime directory if needed
[ ! -d /run/dbus ] && install -m755 -d /run/dbus
exec dbus-daemon --system --nofork --nopidfile
