# etc/sv/agetty-tty5/run
#!/bin/sh
tty=${PWD##*-}
exec chpst -P /sbin/agetty 38400 "${tty}" linux
